# images/lib/consul-server/root/etc/services.d/consul/finish
#!/usr/bin/with-contenv bash

consul leave