# images/lib/glusterfs/root/etc/services.d/glusterfs/run
#!/usr/bin/with-contenv bash

exec glusterd -p /var/run/glusterd.pid -N